# etc/supervisor/conf.d/supervisord-sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -e
user=root
autorestart=true
startretries=3
priority=400